# deploy/supervisor-apps.conf
; /etc/supervisor/conf.d/apps.conf

[program:podcast-digest]
command=/opt/apps/podcast/venv/bin/streamlit run /opt/apps/podcast/podcast_digest/app.py
    --server.port=8501
    --server.headless=true
    --server.address=127.0.0.1
directory=/opt/apps/podcast
autostart=true
autorestart=true
stderr_logfile=/var/log/podcast.err.log
stdout_logfile=/var/log/podcast.out.log
user=ubuntu
environment=HOME="/home/ubuntu"

[program:ativos]
command=/opt/apps/ativos/venv/bin/streamlit run /opt/apps/ativos/app.py
    --server.port=8502
    --server.headless=true
    --server.address=127.0.0.1
directory=/opt/apps/ativos
autostart=true
autorestart=true
stderr_logfile=/var/log/ativos.err.log
stdout_logfile=/var/log/ativos.out.log
user=ubuntu
environment=HOME="/home/ubuntu"

[program:btc-signal]
command=/opt/apps/btc/venv/bin/gunicorn app:server -b 127.0.0.1:8050 --workers 2 --timeout 120
directory=/opt/apps/btc
autostart=true
autorestart=true
stderr_logfile=/var/log/btc.err.log
stdout_logfile=/var/log/btc.out.log
user=ubuntu
